# templates/tomcat_init_script.sh.erb
#!/bin/bash

### BEGIN INIT INFO
# Provides: <%= @service_name %>
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Description: Tomcat init script for <%= @name %>
### END INIT INFO


# chkconfig: 2345 95 20
# description: Tomcat init script for <%= @name %>

export INSTANCE_NAME="<%= @instance_name %>"
export CATALINA_BASE="<%= @instance_dir %>"
export CATALINA_HOME="<%= @catalina_home %>"
LAUNCH_SCRIPT="${CATALINA_HOME}/bin/catalina.sh"
PROCESS_OWNER="<%= @instance_user %>"
WHOAMI=$(whoami)
export CATALINA_PID="<%= @instance_pid %>"
WAIT_TO_DIE=60



function start_instance() {
  $SU_CMD "$LAUNCH_SCRIPT start"
}

function stop_instance() {
  $SU_CMD "$LAUNCH_SCRIPT stop ${WAIT_TO_DIE} -force"
}

function status_instance() {
  if [ -f "$CATALINA_PID" ] ; then
    TC_PROCESS=$(ps "$(cat $CATALINA_PID)")
    PS_STATUS=$?
    if [ "$PS_STATUS" -eq 0 ] ; then
      STATUS=0
      echo "Tomcat instance ${INSTANCE_NAME} running as process $(cat $CATALINA_PID)"
    else
      STATUS=127
      echo "Tomcat instance ${INSTANCE_NAME} should be running as process $(cat $CATALINA_PID) but seems to have died"
    fi
  else
    STATUS=1
    echo "Tomcat instance ${INSTANCE_NAME} is not running"
  fi
  exit $STATUS
}

# check for root (needed to su)
if [ $WHOAMI == "root" ] ; then
  SU_CMD="su $PROCESS_OWNER -c "
elif [ $WHOAMI == "$PROCESS_OWNER" ] ; then
  SU_CMD=""
else
  echo "You must run this script as either root or ${PROCESS_OWNER}"
  exit 1
fi



case "$1" in
  'start')
    start_instance
  ;;
  'stop')
    stop_instance
  ;;
  'restart')
    stop_instance
    start_instance
  ;;
  'status')
    status_instance
  ;;
  *)
    echo "Usage $0 {start|stop|restart|status}"
esac
